# supervisor/qwen36-timing-proxy.conf
[program:qwen36-timing-proxy]
command=/workspace/venvs/sglang/bin/python /opt/supervisor-scripts/qwen36-timing-proxy.py
autostart=true
autorestart=unexpected
exitcodes=0
startsecs=3
stopasgroup=true
killasgroup=true
stopsignal=TERM
stopwaitsecs=20
stdout_logfile=/dev/stdout
redirect_stderr=true
stdout_logfile_maxbytes=0
stdout_logfile_backups=0
