# supervisor/qwen36-public-api.conf
[program:qwen36-public-api]
command=/opt/supervisor-scripts/qwen36-public-api.sh
autostart=true
autorestart=unexpected
exitcodes=0
startsecs=3
stopasgroup=true
killasgroup=true
stopsignal=TERM
stopwaitsecs=20
stdout_logfile=/dev/stdout
redirect_stderr=true
stdout_logfile_maxbytes=0
stdout_logfile_backups=0
