# supervisor/qwen36-nvfp4.conf
[program:qwen36-nvfp4]
environment=PROC_NAME="%(program_name)s"
command=/opt/supervisor-scripts/qwen36-nvfp4-sglang.sh
autostart=true
autorestart=unexpected
exitcodes=0
startsecs=10
stopasgroup=true
killasgroup=true
stopsignal=TERM
stopwaitsecs=180
stdout_logfile=/dev/stdout
redirect_stderr=true
stdout_events_enabled=true
stdout_logfile_maxbytes=0
stdout_logfile_backups=0
